# system/dotfiles-backup.service
[Unit]
Description=Sentry Dotfiles Scheduled Backup Service

[Service]
Type=oneshot
ExecStart=%h/dotfiles/.venv/bin/python %h/dotfiles/backup.py
