# stow/systemd/.config/systemd/user/brainwatch.service
[Unit]
Description=BrainWatch - Personal AI Wellness and Focus System
After=graphical-session.target

[Service]
ExecStart=/home/skc/.local/bin/brainwatch start
Restart=on-failure
RestartSec=5s
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1
Environment=PATH=/home/skc/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
