# stow/systemd/.config/systemd/user/dotfiles-daemon.service
[Unit]
Description=Sentry Dotfiles Change Detection Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/dotfiles/.venv/bin/python %h/dotfiles/daemon.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
